# poweroff.service
[Unit]
Description = Poweroff machine at Midnight

[Service]
Type=oneshot
ExecStart=/usr/sbin/poweroff